# frontdoor/templates/app.service.template
[Unit]
Description=APPNAME_DESCRIPTION
After=network.target tailscaled.service
Wants=tailscaled.service

[Service]
Type=simple
User=APPNAME_USER
WorkingDirectory=APPNAME_INSTALL_DIR
ExecStart=APPNAME_INSTALL_DIR/.venv/bin/uvicorn APPNAME_MODULE:app --host 127.0.0.1 --port APPNAME_INTERNAL_PORT
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
